# Linux/configFiles/enoca.service
[Unit]
Description=EnocaLogger
[Service]
ExecStart=/opt/enoca/enoca.sh
Type=simple
User=alarm
[Install]
WantedBy=multi-user.target
